# devops/var/docker/supervisord/postiz.conf
[program:postiz]
directory=/app/apps/frontend
command=npm run dev
autostart=true
autorestart=true
startsec=3